# supervisor.apps.ini
[program:gunicorn]
command = gunicorn
autostart = true
autorestart = true
startretries = 3
stdout_logfile=logs/gunicorn_out.log
stdout_logfile_maxbytes=10MB
stderr_logfile=logs/gunicorn_err.err
stderr_logfile_maxbytes=10MB

[program:huey]
command = huey_consumer app.scheduler.huey -w 3
autostart = false
autorestart = false
startretries = 3
stdout_logfile=logs/huey_out.log
stdout_logfile_maxbytes=10MB
stderr_logfile=logs/huey_err.err
stderr_logfile_maxbytes=10MB

[program:websockets]
command = python3 app/ws_server
autostart = false
autorestart = false
startretries = 3
